# etc/systemd/system/tv-on.service
[Unit]
Description=Power On TV After Resume
After=suspend.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'sleep 2 && /usr/local/bin/tv-on.sh'

[Install]
WantedBy=suspend.target
